# image/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:cron]
command=/bin/bash -c "cron -f"
autostart=true
autorestart=false